# data/systemd/legion-control-daemon.service
[Unit]
Description=Legion Control hardware daemon
Documentation=https://github.com/daredoole/RatVantage
After=dbus.service multi-user.target

[Service]
Type=dbus
BusName=org.ratvantage.LegionControl1
ExecStart=/usr/libexec/legion-control/legion-control-daemon
Restart=on-failure
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
